# deploy/supervisor-hyuga.ini
[program:hyuga]
process_name=%(program_name)s

directory=/opt/Hyuga
command=/opt/Hyuga/hyuga

startsecs = 5 ; 启动 5 秒后没有异常退出，作为已正常启动
autostart=true
autorestart=true
user=root

redirect_stderr=true
stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20     ; stdout 日志文件备份数
stdout_logfile=/opt/Hyuga/logs/hyuga_stdout.log